# wifi-setup/hotspot.service
[Unit]
Description=Wi-Fi Connectivity Monitor (client first, fallback AP + portal)
After=network-online.target network.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wifi-monitor.sh
Restart=always
RestartSec=5
KillMode=control-group
StandardOutput=journal
StandardError=journal
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
